# patch_management/patch_scanner.service
[Unit]
Description=Patch Management Scanner Service
After=network.target patch_management.service

[Service]
User=patchman
Group=patchman
WorkingDirectory=/opt/patch_management/backend
Environment="PATH=/opt/patch_management/backend/venv/bin"
ExecStart=/opt/patch_management/backend/venv/bin/python patch_scanner.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 