# patch_management/patch_management_frontend.service
[Unit]
Description=Patch Management Frontend Service
After=network.target

[Service]
User=patchman
Group=patchman
WorkingDirectory=/opt/patch_management/frontend
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 